# test/example/simple-php/.webApplication/www-data.crontab
#
# Run once a minute and do something simple
#
* * * * * {APPLICATION_PATH}/bin/cron.php